# rpi-vision-ai-processor/usr/lib/systemd/system/go2rtc.service
[Unit]
Description=go2rtc WebRTC streaming server for Vision AI
Documentation=https://github.com/AlexxIT/go2rtc
After=network.target rpi-vision-ai-processor.service
Wants=rpi-vision-ai-processor.service

[Service]
Type=simple
User=tedge
Group=tedge
WorkingDirectory=/etc/tedge/plugins/vai-plugin/go2rtc
ExecStart=/usr/local/bin/go2rtc -config /etc/tedge/plugins/vai-plugin/go2rtc/go2rtc.yaml
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/tedge/vai-plugin /tmp

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=go2rtc

[Install]
WantedBy=multi-user.target
